# Lab_2/bonus/crontab
* * * * * /media/sf_D_DRIVE/Uni/Programming/Term\ 7/Operating\ Systems/Lab\ 2/bonus/backup-cron.sh source backup 3



